# systemd/user/clash-subscribe.timer
[Unit]
Description=Update clash subscribe config yaml 5 days

[Timer]
OnUnitActiveSec=5d 

[Install]
WantedBy=timers.target
